# supportBox/SBInfo/server/crontab
*/15 * * * * php /var/www/html/supportBox/SBInfo/server/amIOnline.php